# qbittorrent/s6-overlay/etc/services.d/log/run
#!/usr/bin/with-contenv bash

sleep 5
exec s6-setuidgid qbittorrent tail -f /data/data/logs/qbittorrent.log
