# qbittorrent/s6-overlay/etc/services.d/qbittorrent/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid qbittorrent qbittorrent-nox --profile=/home/qbittorrent
